# .ebextensions/supervisord/queue-worker.conf
[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php artisan queue:work sqs --env=production --delay=3 --sleep=3 --tries=3
directory=/var/app/current
autostart=true
autorestart=true
stderr_logfile=/var/app/support/logs/queue-worker-err.log
stdout_logfile=/var/app/support/logs/queue-worker.log
numprocs=1
exitcodes=1
startsecs=0
